# lll.conf
[program:lll]
command=python3 /home/curator2/lll/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/lll/test.err.log
stdout_logfile=/home/curator2/lll/test.out.log